# ha_webdav_mount/rootfs/etc/services.d/webdav/finish
#!/usr/bin/with-contenv bash
echo "[webdav] Service stopped."
